# backend-laravel/deploy/supervisor-reverb.conf
; Keeps the WebSocket server (chat + notifications) running persistently.
; Copy to /etc/supervisor/conf.d/ on the host, adjust user/paths, then:
;   supervisorctl reread && supervisorctl update && supervisorctl start techyfuel-reverb

[program:techyfuel-reverb]
process_name=%(program_name)s
command=php /var/www/techyfuel/backend-laravel/artisan reverb:start
directory=/var/www/techyfuel/backend-laravel
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/techyfuel/backend-laravel/storage/logs/reverb.log
stopasgroup=true
killasgroup=true
