# backend-laravel/deploy/supervisor-queue-worker.conf
; Processes queued jobs — this is what actually delivers broadcasts to
; Reverb (see app/Events/MessagePosted.php, NotificationCreated.php).
; Without this running, chat/notifications are created in the DB but
; never pushed live. Copy to /etc/supervisor/conf.d/ on the host, adjust
; user/paths/numprocs, then:
;   supervisorctl reread && supervisorctl update && supervisorctl start techyfuel-queue-worker:*

[program:techyfuel-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/techyfuel/backend-laravel/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/techyfuel/backend-laravel
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/techyfuel/backend-laravel/storage/logs/queue-worker.log
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
